# devel/etc/catkin/profile.d/99_ksql_airport_setup.sh

# export GAZEBO_MODEL_PATH=/home/xzh/gazebo_build/src/ksql_airport/share/ksql_airport/models:$GAZEBO_MODEL_PATH
export GAZEBO_RESOURCE_PATH=/home/xzh/gazebo_build/src/ksql_airport:$GAZEBO_RESOURCE_PATH
# export GAZEBO_PLUGIN_PATH=/home/xzh/gazebo_build/src/ksql_airport/share/ksql_airport/plugins:$GAZEBO_PLUGIN_PATH
