# devel/etc/catkin/profile.d/99_sand_island_setup.sh

# export GAZEBO_MODEL_PATH=/home/xzh/gazebo_build/src/sand_island/share/sand_island/models:$GAZEBO_MODEL_PATH
export GAZEBO_RESOURCE_PATH=/home/xzh/gazebo_build/src/sand_island:$GAZEBO_RESOURCE_PATH
# export GAZEBO_PLUGIN_PATH=/home/xzh/gazebo_build/src/sand_island/share/sand_island/plugins:$GAZEBO_PLUGIN_PATH
